# server/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:backend]
; Waits for Postgres to accept connections, runs migrations, then starts the
; API. Runs as the unprivileged 'node' user (never as root).
command=sh -c "until pg_isready -h postgres -U ffpro -q; do echo '[backend] waiting for postgres...'; sleep 1; done; node scripts/migrate.js && exec node src/index.js"
directory=/app/server
user=node
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=20
